# rifBuild/LinuxInstall/etc/profile.d/rif.sh
# RIF sh profile
#
export PYTHONPATH=/usr/lib64/python2.6/site-packages/wx-3.0-gtk2
export R_HOME=/usr/local/lib64/R
#
PATH=$PATH:/usr/local/pgsql/bin:/usr/local/pgadmin3/bin
MANPATH=/usr/man:/usr/share/man:/usr/local/man:/usr/local/share/man:/usr/X11R6/man:/usr/local/pgsql/share/man
export PATH MANPATH
#
PGHOST=localhost
export PGHOST
PGDATABASE=sahsuland
export PGDATABASE
#
# Eof
